# conf/systemd.service
[Unit]
Description=Pia Backend
Requires=network.target
Requires=postgresql.service

[Service]
Type=simple
User=root
Group=users
WorkingDirectory=/opt/pia-back
Environment=RAILS_SERVE_STATIC_FILES=true
ExecStart=/bin/bash -c "/opt/pia-back/bin/rails s -e production -b 127.0.0.1 -p 3000"
TimeoutSec=30
RestartSec=15s
Restart=always

[Install]
WantedBy=multi-user.target
